# runit/syncthing/run
#!/bin/sh
exec syncthing serve --no-browser --no-restart 2>&1
